# back/crontab
*/5 * * * * root /usr/local/bin/python /usr/django/app/manage.py send_notifications

